# templates/etc_systemd_system_samba4_sysvol_replication.timer.j2
#
# {{ ansible_managed }}
#
[Unit]
Description=Timer for Samba4 sysvol replication

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
